# systemd/clean-stale-claude-sessions.service
[Unit]
Description=Clean stale Claude sessions from AGS widget

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/clean-stale-claude-sessions.py
