# hosts/mainlaptop/generated/ags_template.service
[Unit]
Description=Launch ags on startup

[Service]
ExecStart=/home/dima/.local/bin/ags run
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitInterval=30
WorkingDirectory=%h/.config/ags
Environment=XDG_CURRENT_DESKTOP=Hyprland
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_HOME=%h/.local/share
Environment=PATH=%h/.local/bin:/usr/bin:%h/go/bin:$PATH
Environment=GTK_THEME=catppuccin-mocha-mauve-standard+default

[Install]
WantedBy=default.target
